# server/horse-app.service
[Unit]
Description=Horse 发生 Node.js Application
Documentation=https://github.com/yourusername/horse-app
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/horse-app/server
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=horse-app

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/horse-app/server/logs

[Install]
WantedBy=multi-user.target
